# profile.d/anaconda3.sh
/ascldap/users/fsabado/dotfiles/anaconda3/anaconda3.sh