# profile.d/pycharm.sh
/ascldap/users/fsabado/dotfiles/pycharm/pycharm.sh