# ddclient/ddclient.cron
# ddclient update
# Place in /etc/cron.d or copy the command using crontab -e for normal users
# Run update every hour
* 1 * * * /usr/sbin/ddclient -verbose -force >> /var/log/ddclient_updates.out 2>&1



